# raspberry_pi_companion/docs/drone-companion.service
[Unit]
Description=Agricultural Drone Companion Computer
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=drone-companion
Group=drone-companion
SupplementaryGroups=video dialout gpio
WorkingDirectory=/opt/drone-companion
Environment="PATH=/opt/drone-companion/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="APP_DATA_DIRECTORY=/var/lib/drone-companion"
Environment="PHOTO_DIRECTORY=/var/lib/drone-companion/photos"
Environment="SPRAY_SESSION_DIRECTORY=/var/lib/drone-companion/spray-sessions"
Environment="MISSION_FILE=/var/lib/drone-companion/missions/mission.json"
Environment="AUDIT_LOG_FILE=/var/lib/drone-companion/audit/events.jsonl"
Environment="PPK_LOG_DIRECTORY=/var/lib/drone-companion/ppk"
Environment="SPRAY_APPLICATION_RECORD_DIRECTORY=/var/lib/drone-companion/application-records"
ExecStart=/opt/drone-companion/venv/bin/python /opt/drone-companion/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=drone-companion

[Install]
WantedBy=multi-user.target
